# systemd/rivendell.service.in
[Unit]
Description=Rivendell Radio Automation System
After=network.target remote-fs.target nss-lookup.target
#Documentation=man:drouterd(8) man:dmap(1) man:dlist(1)

[Service]
Type=simple
ExecStart=@prefix@/sbin/rdservice
PrivateTmp=true
Restart=always

[Install]
WantedBy=multi-user.target
