# debian/postinst
#!/bin/sh
# postinst script for padpoint

case "$1" in
    configure)
      ldconfig
      if test ! -e /etc/padpoint.conf ; then
        cp /usr/share/doc/padpoint/padpoint.conf-sample /etc/padpoint.conf
      fi
      /bin/systemctl daemon-reload
      /bin/systemctl enable padpoint
      /bin/systemctl restart padpoint
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
